# infrastructure/ansible/roles/qrysm-beacon/templates/qrysm-beacon.service.j2
[Unit]
Description=Qrysm beacon chain (QRL consensus client)
Documentation=https://github.com/theQRL/qrysm
After=network-online.target gqrl.service
Wants=network-online.target gqrl.service

[Service]
Type=simple
User={{ qrysm_user }}
Group={{ qrysm_group }}
WorkingDirectory={{ qrysm_beacon_data_dir }}
Restart=on-failure
RestartSec=10
TimeoutStopSec=180

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths={{ qrysm_beacon_data_dir }} /var/log/qrysm

ExecStart={{ qrysm_beacon_binary_path }} \
  --datadir={{ qrysm_beacon_data_dir }} \
  --min-sync-peers={{ qrysm_min_sync_peers }} \
  --genesis-state={{ testnetv2_genesis_path }} \
  --chain-config-file={{ testnetv2_config_path }} \
  --config-file={{ testnetv2_config_path }} \
  --execution-endpoint=http://{{ gqrl_authrpc_addr }}:{{ gqrl_authrpc_port }} \
  --jwt-secret={{ jwt_secret_path }} \
  --accept-terms-of-use \
  --contract-deployment-block=0 \
  --minimum-peers-per-subnet=0 \
  --p2p-static-id \
  --rpc-host={{ qrysm_rpc_host }} --rpc-port={{ qrysm_grpc_port }} \
  --grpc-gateway-host={{ qrysm_grpc_gateway_host }} --grpc-gateway-port={{ qrysm_grpc_gateway_port }} \
  --p2p-tcp-port={{ qrysm_p2p_port }} --p2p-udp-port={{ qrysm_p2p_udp_port }} \
{% if qrysm_metrics_enabled %}
  --monitoring-host={{ qrysm_metrics_addr }} --monitoring-port={{ qrysm_metrics_port }} \
{% endif %}
{% for enr in qrysm_bootstrap_nodes %}
  --bootstrap-node={{ enr }}{% if not loop.last %} \{% endif %}

{% endfor %}

StandardOutput=journal
StandardError=journal
SyslogIdentifier=qrysm-beacon

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
