# infrastructure/ansible/roles/gqrl/templates/gqrl.service.j2
[Unit]
Description=QRL execution client (gqrl)
Documentation=https://github.com/theQRL/go-qrl
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ gqrl_user }}
Group={{ gqrl_group }}
WorkingDirectory={{ gqrl_data_dir }}
Restart=on-failure
RestartSec=10
TimeoutStopSec=180

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths={{ gqrl_data_dir }} /var/log/gqrl

ExecStart={{ gqrl_binary_path }} \
  --networkid {{ gqrl_chain_id }} \
  --datadir={{ gqrl_data_dir }} \
  --nat=extip:{{ gqrl_nat_extip }} \
  --port {{ gqrl_port }} \
  --syncmode={{ gqrl_syncmode }} \
  --http \
  --http.addr={{ gqrl_http_addr }} \
  --http.port={{ gqrl_http_port }} \
  --http.api={{ gqrl_http_api }} \
  --authrpc.addr={{ gqrl_authrpc_addr }} \
  --authrpc.port={{ gqrl_authrpc_port }} \
  --authrpc.jwtsecret={{ jwt_secret_path }} \
{% if gqrl_metrics_enabled %}
  --metrics --metrics.addr={{ gqrl_metrics_addr }} --metrics.port={{ gqrl_metrics_port }}
{% endif %}

StandardOutput=journal
StandardError=journal
SyslogIdentifier=gqrl

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
