# infrastructure/ansible/roles/qrysm-validator/templates/qrysm-validator.service.j2
[Unit]
Description=Qrysm validator client (QRL staking)
Documentation=https://github.com/theQRL/qrysm
After=network-online.target qrysm-beacon.service
Wants=network-online.target qrysm-beacon.service

[Service]
Type=simple
User={{ qrysm_user }}
Group={{ qrysm_group }}
WorkingDirectory={{ qrysm_validator_data_dir }}
Restart=on-failure
RestartSec=10
TimeoutStopSec=180

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths={{ qrysm_validator_data_dir }} /var/log/qrysm

# CRITICAL: Pre-start check to prevent running on backup nodes
ExecStartPre=/usr/local/bin/qrysm-validator-precheck

ExecStart={{ qrysm_validator_binary_path }} \
  --datadir={{ qrysm_validator_data_dir }} \
  --wallet-dir={{ validator_wallet_dir }} \
{% if validator_wallet_password_file is defined and validator_wallet_password_file != "" %}
  --wallet-password-file={{ validator_wallet_password_file }} \
{% endif %}
  --chain-config-file={{ testnetv2_config_path }} \
  --config-file={{ testnetv2_config_path }} \
  --beacon-rpc-provider={{ qrysm_beacon_rpc_provider }} \
  --beacon-rpc-gateway-provider={{ qrysm_beacon_rpc_gateway_provider }} \
  --accept-terms-of-use \
  --suggested-fee-recipient={{ fee_recipient_address }} \
{% if qrysm_validator_metrics_enabled %}
  --monitoring-host={{ qrysm_validator_metrics_addr }} --monitoring-port={{ qrysm_validator_metrics_port }}
{% endif %}

StandardOutput=journal
StandardError=journal
SyslogIdentifier=qrysm-validator

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
